# supervisord.ini
[program:sgviewer]
command=bash -c '. /Volumes/VFX/production/key_tools/key_base/systems/environment/environment.sh; make serve'
directory=/www/sgviewer
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
environment=PORT=80

